# fsdeploy/contrib/openrc/fsdeploy.init
#!/sbin/openrc-run
# Permissions: chmod +x

name="fsdeploy"
description="fsdeploy service"

command="/usr/bin/fsdeploy"
command_args="--daemon"
command_background=true
pidfile="/var/run/${RC_SVCNAME}.pid"

depend() {
    need localmount
    after bootmisc
}
